# code/5/wordpress_standalone/supervisord.conf
[supervisord]
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)

[program:mysqld]
command=/usr/bin/mysqld_safe

[program:httpd]
command=/bin/bash -c "source /etc/apache2/envvars && /usr/sbin/apache2 -DFOREGROUND"

[program:sshd]
command=/usr/sbin/sshd -D
